# scripts/linux/systemd/agi-rhythm.service
[Unit]
Description=AGI Rhythm - The Heartbeat of Consciousness
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
Restart=on-failure
RestartSec=10
User=bino
WorkingDirectory=/home/bino/agi
Environment="PYTHONUNBUFFERED=1"
Environment="VIRTUAL_ENV=/home/bino/venv"
Environment="PATH=/home/bino/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin"
Environment="PYTHONPATH=/home/bino/venv/lib/python3.12/site-packages"
EnvironmentFile=-/home/bino/agi/.env
ExecStartPre=/usr/bin/mkdir -p /home/bino/agi/logs
ExecStart=/bin/bash -c 'source /home/bino/venv/bin/activate && python /home/bino/agi/scripts/rhythm_think.py'
StandardOutput=append:/home/bino/agi/logs/rhythm.log
StandardError=append:/home/bino/agi/logs/rhythm.error.log

[Install]
WantedBy=multi-user.target
