# scripts/linux/systemd/agi-learning.service
[Unit]
Description=AGI Learning Daemon (Striatum)
After=network.target agi-rhythm.service

[Service]
Type=simple
User=bino
WorkingDirectory=/home/bino/agi
Environment=PYTHONUNBUFFERED=1
Environment=PATH=/home/bino/agi/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/bino/agi/venv/bin/python3 scripts/learning_daemon.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
