# scripts/linux/systemd/agi-rhythm-user.service
[Unit]
Description=AGI Rhythm - The Heartbeat of Consciousness
After=network.target

[Service]
Type=simple
Restart=on-failure
RestartSec=10
WorkingDirectory=/home/bino/agi
Environment="PYTHONUNBUFFERED=1"
ExecStart=/bin/bash -c 'source /home/bino/venv/bin/activate && python /home/bino/agi/scripts/rhythm_think.py'
StandardOutput=append:/home/bino/agi/logs/rhythm.log
StandardError=append:/home/bino/agi/logs/rhythm.error.log

[Install]
WantedBy=default.target
