# systemd/agi-axiom-sync.service
[Unit]
Description=AGI Axiom Synchronization Daemon
Documentation=https://github.com/ruafieldphase/agi
After=network.target

[Service]
Type=simple
ExecStart=/home/bino/.agi_venv/bin/python3 /home/bino/agi/scripts/axiom_sync_daemon.py
WorkingDirectory=/home/bino/agi
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agi-axiom-sync

[Install]
WantedBy=multi-user.target
